# fan-control.service
[Unit]
Description=UCG-Max/Fibre Adaptive Fan Controller
After=basic.target
Documentation=https://github.com/iceteaSA/ucg-max-fan-control

[Service]
Type=simple
User=root
WorkingDirectory=/data/fan-control
ExecStart=/data/fan-control/fan-control.sh
Restart=always
RestartSec=5
# 24h window for log rate limiting
LogRateLimitIntervalSec=86400
# Allow 10k logs/day
LogRateLimitBurst=10000

[Install]
WantedBy=multi-user.target
