# fan-control.service
[Unit]
Description=UCG-Max Adaptive Fan Controller
After=basic.target
Documentation=https://github.com/iceteaSA/ucg-max-fan-control

[Service]
Type=simple
User=root
WorkingDirectory=/data/fan-control
ExecStart=/data/fan-control/fan-control.sh
Restart=always
RestartSec=5
LogRateLimitIntervalSec=86400  # 24h window
LogRateLimitBurst=10000       # Allow 10k logs/day

[Install]
WantedBy=multi-user.target
